# touchscreen/udev-rules/80-touchscreen.rules
SUBSYSTEM=="drm", ACTION=="change", DRIVERS=="?*", KERNEL=="card*", RUN+="/usr/local/bin/conf-touchscreen.sh"
